# etc/service/boringtun/run
#!/usr/bin/env bash

set -euo pipefail

exec /usr/local/sbin/boringtun --foreground wgnet0
