# etc/service/sshd/run
#!/usr/bin/env bash

exec /usr/sbin/sshd -E - -D
